# novosga/config/crontab
# NovoSGA Crontab
# Run every minute
* * * * * cd /var/www/html && php bin/console novosga:cleanup || true
